# script/cron/crontab
30  *   *   *   *   checkupdates | awk 'END{print NR}' > /home/$USER/updates.txt
